# services/picam-freeze-watchdog.cron
*/3 * * * * root /home/pi/picam/watchdog.sh
