# services/_deprecated/picam-watchdog.timer
[Unit]
Description=PiCam Watchdog Timer - Runs health checks every 5 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
